# etc/init.d/power-profiles-daemon
#!/sbin/openrc-run

description="Manages power profiles for modern desktops"
command="/usr/libexec/power-profiles-daemon"
supervisor=supervise-daemon

depend() {
	need dbus
	provide power-profiles-daemon
}

start_pre() {
	checkpath -d -o root:root -m 0755 /var/lib/power-profiles-daemon
}
