# etc/init.d/power-profiles-daemon-openrc
#!/sbin/openrc-run

description="Power Profiles Daemon"

command="/usr/libexec/power-profiles-daemon"
pidfile="/run/power-profiles-daemon.pid"

depend() {
    need dbus
    after multi-user display-manager
    provide power-management
    conflict tuned tlp auto-cpufreq system76-power
}

start() {
    ebegin "Starting Power Profiles Daemon"

    checkpath --directory --owner root:root --mode 0755 /var/lib/power-profiles-daemon
    checkpath --file --owner root:root --mode 0644 /run/power-profiles-daemon.pid

    start-stop-daemon --start --exec "${command}" --make-pidfile --pidfile "${pidfile}" --background
    eend $?
}

stop() {
    ebegin "Stopping Power Profiles Daemon"
    start-stop-daemon --stop --exec "${command}" --pidfile "${pidfile}"
    eend $?
}

restart() {
    ebegin "Restarting Power Profiles Daemon"
    stop
    start
    eend $?
}
