# etc/init.d/power-profiles-daemon-openrc
#!/sbin/openrc-run

description="Power Profiles Daemon"

command="/usr/libexec/power-profiles-daemon"
pidfile="/run/power-profiles-daemon.pid"

depend() {
    need dbus
    after sddm
}

start() {
    ebegin "Starting Power Profiles Daemon"

    checkpath --directory --owner root:root --mode 0755 /var/lib/power-profiles-daemon
    checkpath --directory --owner root:root --mode 0755 /run
    checkpath --file --owner root:root --mode 0644 "${pidfile}"

    start-stop-daemon --start --exec "${command}" --make-pidfile --pidfile "${pidfile}" --background
    eend $?
}

stop() {
    ebegin "Stopping Power Profiles Daemon"
    start-stop-daemon --stop --pidfile "${pidfile}"
    eend $?
}

restart() {
    ebegin "Restarting Power Profiles Daemon"
    stop
    sleep 1
    start
    eend $?
}
