# debian/postinst
#!/usr/bin/env bash

# -- Add service to runlevel.

rc-update add power-profiles-daemon-openrc default

# -- Update services list.

rc-update -u
